# scripts/nxr-admin.service
[Unit]
Description=NXR Admin Flask App
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/nxr_website
EnvironmentFile=-/root/nxr_website/.env.admin
Environment=PORT=8081
Environment=FLASK_DEBUG=0
Environment=PYTHONUNBUFFERED=1
ExecStart=/usr/bin/python3 /root/nxr_website/nxr_admin/app_updated.py
Restart=always
RestartSec=3
StandardOutput=append:/root/nxr_website/nxr_admin/output.log
StandardError=append:/root/nxr_website/nxr_admin/output.log

[Install]
WantedBy=multi-user.target
